# deploy/systemd/quadball-timer.service
[Unit]
Description=Quadball Timer Bun server
After=network.target

[Service]
Type=simple
User=quadball-timer
Group=quadball-timer
WorkingDirectory=/srv/quadball-timer/current
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3000
ExecStart=/srv/quadball-timer/current/quadball-timer
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
UMask=027

[Install]
WantedBy=multi-user.target
